# systemd/lib/systemd/system/refresh-btcmap-data.service
[Unit]
Description=Refresh BTC Map data

[Service]
ExecStart=/usr/local/bin/refresh-btcmap-data
